# lanmp/nginx/nginx/files
etc/
etc/default/
etc/default/nginx
etc/logrotate.d/
etc/logrotate.d/nginx
etc/nginx/
etc/nginx/html
etc/nginx/logs
etc/nginx/conf/
etc/nginx/conf/fastcgi.conf
etc/nginx/conf/fastcgi.conf.default
etc/nginx/conf/fastcgi_params
etc/nginx/conf/fastcgi_params.default
etc/nginx/conf/koi-utf
etc/nginx/conf/koi-win
etc/nginx/conf/mime.types
etc/nginx/conf/mime.types.default
etc/nginx/conf/nginx.conf
etc/nginx/conf/nginx.conf.default
etc/nginx/conf/scgi_params
etc/nginx/conf/scgi_params.default
etc/nginx/conf/uwsgi_params
etc/nginx/conf/uwsgi_params.default
etc/nginx/conf/win-utf
etc/init.d/
etc/init.d/nginx
srv/
srv/http/
srv/http/nginx/
srv/http/nginx/50x.html
usr/
usr/sbin/
usr/sbin/nginx
usr/share/
usr/share/licenses/
usr/share/licenses/nginx/
usr/share/licenses/nginx/LICENSE
var/
var/log/
var/log/nginx/
var/spool/
var/spool/nginx/


